# deployments/dev/systemd/container-geohod-postgres-dev.service
[Unit]
Description=Geohod Postgres container (dev)
After=network.target
Wants=geohod-dev-pod.service

[Service]
Restart=on-failure
RestartSec=5
ExecStartPre=/usr/bin/podman pull docker.io/library/postgres:17-alpine
ExecStart=/usr/bin/podman run --rm --name geohod-postgres-dev \
  --pod geohod-dev-pod \
  -e POSTGRES_DB=${POSTGRES_DB} \
  -e POSTGRES_USER=${POSTGRES_USER} \
  -e POSTGRES_PASSWORD_FILE=/etc/geohod/pg-password \
  -v /var/lib/geohod/postgres:/var/lib/postgresql/data \
  docker.io/library/postgres:17-alpine
ExecStop=/usr/bin/podman stop geohod-postgres-dev || true
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target